# deploy/host/rclone-rcd.service
# rclone remote-control daemon for tvh-controller uploads.
# Install on each tvheadend host:
#   1. install the rclone binary (https://rclone.org/install/)
#   2. place rclone.conf with the Google Drive remote (same one used manually)
#   3. copy this file to /etc/systemd/system/rclone-rcd.service, edit User
#      and paths
#   4. systemctl daemon-reload && systemctl enable --now rclone-rcd
#
# SECURITY: --rc-no-auth means this daemon accepts ANY rc request from ANY
# host that can reach --rc-addr, with no credentials and no authorization
# checks (it can trigger arbitrary file copies/moves/deletes via the rc
# API). This is only safe because --rc-addr is bound to an interface that is
# reachable exclusively over the zone VPN, and the controller is the only
# client expected to reach it — the daemon itself provides NO defense in
# depth. Do not bind --rc-addr to a publicly routable interface, do not
# expose port 5572 through any firewall/NAT, and do not remove --rc-no-auth
# without also adding `user`/`pass` to the corresponding instance's rclone
# config in tvh-controller's config.yaml. --transfers 1 keeps host I/O and
# Drive quota impact low.
[Unit]
Description=rclone remote control daemon (tvh-controller uploads)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hts
ExecStart=/usr/bin/rclone rcd \
  --rc-addr 0.0.0.0:5572 \
  --rc-no-auth \
  --config /home/hts/.config/rclone/rclone.conf \
  --transfers 1 \
  --bwlimit 5M \
  --log-level NOTICE
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
